# deploy/healthcheck.service
[Unit]
Description=Infrastructure Health Check Service
After=network.target

[Service]
Type=simple
ExecStart=%h/life/repos/healthcheck/healthcheck -config %h/life/repos/healthcheck/config.yaml
Restart=on-failure
RestartSec=10
Environment=HOME=%h

[Install]
WantedBy=default.target
